# gandi-update.timer
[Unit]
Description=Update IP address records of domains on Gandi

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

